# Minecraft Dev Docker Client/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:fabric-setup]
command=/home/minecraft/install_fabric.sh
user=minecraft
autorestart=false
startsecs=0
priority=100

[program:forge-setup]
command=/home/minecraft/install_forge.sh
user=minecraft
autorestart=false
startsecs=0
priority=101

[program:copy-mods]
command=/home/minecraft/copy_latest_mod.sh
user=minecraft
autorestart=false
startsecs=0
priority=102

[program:fabric-server]
command=java -Xmx2G -Xms1G -jar fabric-server-mc.1.20.1-loader.0.15.11-launcher.1.0.1.jar nogui
directory=/home/minecraft/fabric-server
user=minecraft
autorestart=true
startsecs=10
priority=200

[program:forge-server]
command=java -Xmx2G -Xms1G -jar forge-1.20.1-47.3.0.jar nogui
directory=/home/minecraft/forge-server
user=minecraft
autorestart=true
startsecs=10
priority=201
